# orin/recomo_service_control.service
[Unit]
Description=Recomo Service Control API
After=network.target

[Service]
Type=simple
User=cr
WorkingDirectory=/home/nvidia/videoFromAndroid/orin
Environment="SERVICE_CONTROL_PIN=1234"
Environment="SERVICE_CONTROL_PORT=8083"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 -m uvicorn service_control_api:app --host 0.0.0.0 --port 8083
Restart=always
RestartSec=5
StandardOutput=append:/home/nvidia/videoFromAndroid/orin/service_control_api.log
StandardError=append:/home/nvidia/videoFromAndroid/orin/service_control_api.log

[Install]
WantedBy=multi-user.target
